# env/supervisord.d/binlog.ini
[program:binlog]
directory = /home/ds/oblogproxy
command = bash -c "sleep 1 && export LD_LIBRARY_PATH=/home/ds/oblogproxy/deps/lib && echo ${LD_LIBRARY_PATH} && ulimit -c unlimited && ulimit -d unlimited && ulimit -f unlimited && ulimit -m unlimited && ulimit -t unlimited && ulimit -v unlimited && ulimit -x unlimited && ulimit -n 102400 && ulimit -q 819200 && ulimit -i 2062651 && ulimit -u 65535 && ulimit -s 8192 && ulimit -l 64 && ulimit -r 0 && cd /home/ds/oblogproxy && ldd ./bin/logproxy && chmod u+x ./bin/logproxy && ./bin/logproxy -f ./conf/conf.json &> log/out.log"
autostart = true
autorestart = true
startsecs = 0
startretries = 9999
stopsignal = QUIT
stderr_logfile = /home/ds/log/oblogproxy_stderr.log
stdout_logfile = /home/ds/log/oblogproxy_stdout.log