# env/supervisord.d/oblogproxy.ini
[program:oblogproxy]
directory = /home/ds/oblogproxy
command = bash -c "sleep 3 && export LD_LIBRARY_PATH=/home/ds/oblogproxy/deps/lib && ulimit -c unlimited && cd /home/ds/oblogproxy && ldd ./bin/logproxy && chmod u+x ./bin/logproxy && ./bin/logproxy -f ./conf/conf.json &> log/out.log"
autostart = true
autorestart = true
startsecs = 0
startretries = 9999
stopsignal = QUIT
stderr_logfile = /home/ds/log/oblogproxy_stderr.log
stdout_logfile = /home/ds/log/oblogproxy_stdout.log